# backend/supervisor.conf
[supervisord]
logfile=/dev/stdout 
logfile_maxbytes=0  
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root
environment=PYTHONUNBUFFERED=1,FLASK_APP="hallway-display",GUNICORN_CMD_ARGS="--bind=[::]:4999 --workers=2",SL_API_KEY="%(ENV_SL_API_KEY)s",OPENWEATHER_API_KEY="%(ENV_OPENWEATHER_API_KEY)s",API_AUTH_TOKEN="%(ENV_API_AUTH_TOKEN)s",HOME_COORDINATES="%(ENV_HOME_COORDINATES)s"

[unix_http_server]
file=/tmp/supervisor.sock

[program:server]
command=./start_server.sh
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true